# simulator/4.4-forward/etc/service/cloudstack_simulator/run
#!/bin/bash

cd /root/cloudstack
exec mvn -Dsimulator -pl client jetty:run > /var/log/cloudstack_simulator.log 2>&1
